# overlay-includes.chroot/etc/systemd/system/security-init.service
[Unit]
Description=Security Initialization and Hardening Service
DefaultDependencies=no
After=local-fs.target systemd-sysctl.service
Before=multi-user.target display-manager.service
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/security-init.sh

# Security hardening
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_SYS_NICE CAP_NET_ADMIN CAP_DAC_OVERRIDE CAP_SETUID CAP_SETGID CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target

